# rootfs/etc/s6/registry/run
#!/bin/bash

if [ -f ./setup ]
then
  ./setup
fi

STARTCMD="/usr/bin/registry /etc/registry/config.yml"
exec ${STARTCMD}
